# build/Immortalwrt/diy/package/base-files/files/etc/uci-defaults/99-my-custom-config
#!/bin/sh
uci set dhcp.lan.leasetime='2m'

uci add dhcp host # =cfg05fe63
uci add_list dhcp.@host[-1].mac='6A:CE:3B:D4:CB:8B'
uci set dhcp.@host[-1].ip='192.168.10.158'
uci add dhcp host # =cfg06fe63
uci add_list dhcp.@host[-1].mac='30:9C:23:E1:E3:72'
uci set dhcp.@host[-1].ip='192.168.10.107'
uci add dhcp host # =cfg07fe63
uci add_list dhcp.@host[-1].mac='7C:2B:E1:13:6E:83'
uci set dhcp.@host[-1].ip='192.168.10.180'

uci set autoreboot.cfg01f8be.enable='1'
uci set autoreboot.cfg01f8be.week='2'
uci set autoreboot.cfg01f8be.hour='5'

uci del zerotier.sample_config.join
uci add_list zerotier.sample_config.join='d3ecf5726da3eeac'
uci set zerotier.sample_config.nat='1'

uci set vlmcsd.config.enabled='1'
uci set vlmcsd.config.auto_activate='1'
uci set vlmcsd.config.internet_access='1'

uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6=node
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.label='233boy-hysteria2-67.216.192.27'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.type='hysteria2'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.address='67.216.192.27'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.port='2022'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.password='9462b0f3-19e8-4ea3-bb7b-4cd7a5e958a4'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.tls='1'
uci set homeproxy.0a5c7eb629c19bdc98ad795e159758c6.tls_insecure='1'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba=node
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.label='233boy-reality-67.216.192.27'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.type='vless'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.address='67.216.192.27'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.port='64841'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.uuid='561b0cb6-ed7a-4ae6-a7b3-b088094f939e'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.tls='1'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.tls_sni='www.paypal.com'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.tls_reality='1'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.tls_reality_public_key='exEUfLoUxz_fMd-vNtu8SW-0YMI4XJCu0MU3xSTP1k8'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.tls_utls='chrome'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.vless_flow='xtls-rprx-vision'
uci set homeproxy.867c7a79333ea464fbbde4a7cde7a0ba.packet_encoding='xudp'
uci set homeproxy.b5916c7711968726d341bac5660944d6=node
uci set homeproxy.b5916c7711968726d341bac5660944d6.label='233boy-tcp-67.216.192.27'
uci set homeproxy.b5916c7711968726d341bac5660944d6.type='vmess'
uci set homeproxy.b5916c7711968726d341bac5660944d6.address='67.216.192.27'
uci set homeproxy.b5916c7711968726d341bac5660944d6.port='6501'
uci set homeproxy.b5916c7711968726d341bac5660944d6.uuid='6eb22e1d-cd61-4960-b418-21e9fb1ccf72'
uci set homeproxy.b5916c7711968726d341bac5660944d6.vmess_alterid='0'
uci set homeproxy.b5916c7711968726d341bac5660944d6.vmess_encrypt='auto'
uci set homeproxy.b5916c7711968726d341bac5660944d6.tls='0'
uci set homeproxy.b5916c7711968726d341bac5660944d6.packet_encoding='xudp'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa=node
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.label='233boy-hysteria2-66.103.198.190'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.type='hysteria2'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.address='66.103.198.190'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.port='5505'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.password='14f73ac8-3f9b-4884-a5e5-d248a4ec8df7'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.tls='1'
uci set homeproxy.cc61cdaf5185c3ed11295f19173787aa.tls_insecure='1'
uci set homeproxy.0147b9839bf8212ad73c632321b94946=node
uci set homeproxy.0147b9839bf8212ad73c632321b94946.label='233boy-tcp-66.103.198.190'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.type='vmess'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.address='66.103.198.190'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.port='8086'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.uuid='7222f228-b7fb-4755-9ffc-73105adda445'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.vmess_alterid='0'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.vmess_encrypt='auto'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.tls='0'
uci set homeproxy.0147b9839bf8212ad73c632321b94946.packet_encoding='xudp'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299=node
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.label='233boy-tcp-104.168.13.17'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.type='vmess'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.address='104.168.13.17'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.port='33034'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.uuid='1c41ab36-d5b9-4d32-b55f-ac50536df38a'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.vmess_alterid='0'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.vmess_encrypt='auto'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.tls='0'
uci set homeproxy.d5a2e4be03d3fd07fc14cb5744465299.packet_encoding='xudp'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a=node
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.label='233boy-hysteria2-104.168.13.17'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.type='hysteria2'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.address='104.168.13.17'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.port='33033'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.password='df7220dd-e3a6-4432-a6dd-557eb6a07fa5'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.tls='1'
uci set homeproxy.22e8c8887a6e29889ad8a440f648f29a.tls_insecure='1'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a=node
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.label='233boy-reality-104.168.13.17'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.type='vless'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.address='104.168.13.17'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.port='30237'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.uuid='d53e672e-1b57-4697-a267-beb4ba36356e'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.tls='1'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.tls_sni='www.ebay.com'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.tls_reality='1'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.tls_reality_public_key='R1z_NDAs4Lm61ZQbttqDAF6ZXDfKspkl6rVCpRoCMH8'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.tls_utls='chrome'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.vless_flow='xtls-rprx-vision'
uci set homeproxy.be84ab866b4948d3c1b1573d93a83b6a.packet_encoding='xudp'
uci del homeproxy.dns.disable_cache
uci del homeproxy.dns.disable_cache_expire
uci set homeproxy.config.main_node='urltest'
uci add_list homeproxy.config.main_urltest_nodes='0a5c7eb629c19bdc98ad795e159758c6'
uci add_list homeproxy.config.main_urltest_nodes='867c7a79333ea464fbbde4a7cde7a0ba'
uci add_list homeproxy.config.main_urltest_nodes='b5916c7711968726d341bac5660944d6'
uci add_list homeproxy.config.main_urltest_nodes='cc61cdaf5185c3ed11295f19173787aa'
uci add_list homeproxy.config.main_urltest_nodes='0147b9839bf8212ad73c632321b94946'
uci add_list homeproxy.config.main_urltest_nodes='d5a2e4be03d3fd07fc14cb5744465299'
uci add_list homeproxy.config.main_urltest_nodes='22e8c8887a6e29889ad8a440f648f29a'
uci add_list homeproxy.config.main_urltest_nodes='be84ab866b4948d3c1b1573d93a83b6a'
uci set homeproxy.config.china_dns_server='223.5.5.5'
uci set homeproxy.routing.udp_timeout='300'
uci set homeproxy.routing.bypass_cn_traffic='0'
uci set homeproxy.dns.default_server='default-dns'














exit 0  # 必须添加，确保执行成功后被删除
